# server/squarezd.service.in
[Unit]
Description=Squarez daemon
After=network.target

[Service]
User=root
Type=simple
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/squarezd

[Install]
WantedBy=multi-user.target