# server/squarezd.service.in
[Unit]
Description=Squarez daemon
After=network.target

[Service]
User=squarez
Type=simple
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/squarezd
Restart=on-failure

[Install]
WantedBy=multi-user.target