# server/squarezd.socket
[Unit]
Description=Squarez daemon socket

[Socket]
ListenStream=3000
Accept=no

[Install]
WantedBy=sockets.target